# utils/smtp_check.crontab
0 18 * * * /usr/bin/bash /home/astun/docker-geonetwork/utils/smtp_check.sh
